# NoticeWorkerService/script/happyshop.service
[Unit]
Description=Run my script every 5 seconds

[Service]
Type=oneshot
ExecStart=/home/service/happyshop/script.sh